# packaging/service-adaptor.service
[Unit]
Description=Service Adaptor daemon
Requires=dbus.socket
After=dlog_logger.service tizen-runtime.target

[Service]
User=system
Group=system
SmackProcessLabel=service-adaptor
Type=dbus
BusName=com.samsung.svcadaptor
ExecStart=/usr/bin/service-adaptor
Restart=none
RestartSec=0
KillMode=none

[Install]
WantedBy=multi-user.target
